# beeminder/jrnl-beeminder-watcher.path
[Path]
Unit=jrnl-beeminder-watcher.service
PathChanged=/home/dominik/Sync/jrnl

[Install]
WantedBy=multi-user.target
